# etc/sv/kafka/run
#!/bin/sh

echo " starting kafka"
exec /opt/kafka/bin/kafka-server-start.sh /etc/kafka/conf/server.properties